# overlay/etc/supervisor.d/accounts-daemon.ini

[program:accounts-daemon]
priority=12
autostart=true
autorestart=true
directory=/
user=root
command=/usr/libexec/accounts-daemon
environment=HOME="/root",USER="root"
stopsignal=TERM
stdout_logfile=/var/log/supervisor/accounts-daemon.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/accounts-daemon.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
