# overlay/etc/supervisor.d/polkit.ini

[program:polkit]
priority=11
autostart=true
autorestart=true
directory=/
user=root
command=/usr/lib/polkit-1/polkitd --no-debug
environment=HOME="/root",USER="root"
stopsignal=TERM
stdout_logfile=/var/log/supervisor/polkit.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/polkit.err.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=3
